# systemd-scripts/plstalk.service
[Unit]
Description=Ftrace syscalls

[Service]
Type=forking
ExecStart=/usr/bin/myscript

[Install]
WantedBy=multi-user.target
