# systemd-scripts/timer.timer
[Unit]
Description=timer for ftrace timer

[Timer]
OnBootSec=2h

[Install]
WantedBy=timers.target
