# systemd-scripts/timer.service
[Unit]
Description=dumps ftrace

[Service]
Type=oneshot
ExecStart=/usr/bin/timescript
